# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:termix]
command=node dist/backend/starter.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/termix_stdout.log
stderr_logfile=/var/log/supervisor/termix_stderr.log
stopasgroup=true
killasgroup=true

[program:r2-sync]
command=/bin/bash /sync.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/r2sync_stdout.log
stderr_logfile=/var/log/supervisor/r2sync_stderr.log
stopasgroup=true
killasgroup=true
